# logstash.conf
description "Logstash log collector"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

limit nofile 8192 8192

exec /usr/bin/logstash system_agent

